# wkrt-fm.service
[Unit]
Description=WKRT-FM 104.7 — Retro Radio Engine
After=network.target icecast2.service
Wants=icecast2.service

[Service]
Type=simple
User=roxanne
WorkingDirectory=/home/roxanne/wkrt
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/home/roxanne/wkrt/.env
ExecStart=/usr/bin/python3 /home/roxanne/wkrt/main.py
Restart=on-failure
RestartSec=15

[Install]
WantedBy=multi-user.target
