# wkrt-ingest.path
[Unit]
Description=Watch for new WKRT music drops in new_music/

[Path]
PathChanged=/home/roxanne/wkrt/new_music
Unit=wkrt-ingest.service

[Install]
WantedBy=multi-user.target
